# roles/github_runner/templates/cleanup-docker.service.j2
# Managed by Ansible — do not edit manually.

[Unit]
Description=Prune stale Docker resources

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cleanup-docker.sh
